# supervisor.conf
[program:sshd]
command=/usr/bin/mysqld_safe
